# infrastructure/tactical-rmm/tactical-base/supervisor.conf
[supervisord]
nodaemon=true
logfile=${TACTICAL_DIR}/temp/supervisord.log
pidfile=${TACTICAL_DIR}/temp/supervisord.pid

[include]
files = ${TACTICAL_DIR}/supervisor/*.conf

[program:nats-server]
command=nats-server --config ${NATS_CONFIG}
stdout_logfile=${TACTICAL_DIR}/logs/nats-server.log
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:agent-config-watcher]
command=/bin/bash ${TACTICAL_DIR}/supervisor/config_watcher.sh
startsecs=1
autorestart=true
startretries=1
stdout_logfile=${TACTICAL_DIR}/logs/config-watcher.log
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:nats-api]
command=/bin/bash -c "/usr/local/bin/nats-api -config ${NATS_API_CONFIG}"
stdout_logfile=${TACTICAL_DIR}/logs/nats-api.log
stdout_logfile_maxbytes=0
redirect_stderr=true
