# infrastructure/tactical-rmm/tactical-base/supervisor.conf
[unix_http_server]
file=/var/run/supervisor/supervisor.sock
chmod=0660
chown=tactical:tactical

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisor/supervisord.pid
user=tactical

[supervisorctl]
serverurl=unix:///var/run/supervisor/supervisor.sock

[include]
files = ${TACTICAL_DIR}/supervisor/*.conf

[program:nats-server]
command=/usr/local/bin/nats-server -c ${NATS_CONFIG}
user=tactical
autostart=true
autorestart=true
startretries=5
startsecs=5
stopwaitsecs=10
stdout_logfile=/var/log/supervisor/nats-server.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
redirect_stderr=true
environment=NATS_HOME="${TACTICAL_DIR}",PATH="/usr/local/bin:/usr/bin:/bin"

[program:agent-listener]
command=${TACTICAL_DIR}/supervisor/agent_listener.sh
user=tactical
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=10
stdout_logfile=/var/log/supervisor/agent-listener.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
redirect_stderr=true
environment=PYTHONPATH="${TACTICAL_DIR}/api",DJANGO_SETTINGS_MODULE="tacticalrmm.settings"

[program:nats-api]
command=/usr/local/bin/nats-api -c ${NATS_API_CONFIG}
user=tactical
autostart=true
autorestart=true
startretries=5
startsecs=5
stopwaitsecs=10
stdout_logfile=/var/log/supervisor/nats-api.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
redirect_stderr=true
