# infrastructure/tactical-rmm/tactical-base/supervisor.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
user=tactical

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[include]
files = ${TACTICAL_DIR}/supervisor/*.conf

[program:nats-server]
command=/usr/local/bin/nats-server -c ${NATS_CONFIG}
user=tactical
autostart=true
autorestart=true
startretries=5
startsecs=5
stopwaitsecs=10
stderr_logfile=/var/log/supervisor/nats-server.err.log
stdout_logfile=/var/log/supervisor/nats-server.out.log
redirect_stderr=true
environment=NATS_HOME="${TACTICAL_DIR}",PATH="/usr/local/bin:/usr/bin:/bin"

[program:agent-listener]
command=${TACTICAL_DIR}/supervisor/agent_listener.sh
user=tactical
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=10
stderr_logfile=/var/log/supervisor/agent-listener.err.log
stdout_logfile=/var/log/supervisor/agent-listener.out.log
redirect_stderr=true
environment=PYTHONPATH="${TACTICAL_DIR}/api",DJANGO_SETTINGS_MODULE="tacticalrmm.settings"

[program:nats-api]
command=/usr/local/bin/nats-api -c ${NATS_API_CONFIG}
user=tactical
autostart=true
autorestart=true
startretries=5
startsecs=5
stopwaitsecs=10
stderr_logfile=/var/log/supervisor/nats-api.err.log
stdout_logfile=/var/log/supervisor/nats-api.out.log
redirect_stderr=true
