# infrastructure/tactical-rmm/tactical-nats/supervisord.conf
[supervisord]
nodaemon=true
user=tactical
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:nats]
command=nats-server -c /opt/tactical/nats/nats-rmm.conf
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/nats.err.log
stdout_logfile=/var/log/supervisor/nats.out.log

[program:config_watcher]
command=/usr/local/bin/config_watcher.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/config_watcher.err.log
stdout_logfile=/var/log/supervisor/config_watcher.out.log 